# files/horizon/etc/services.d/horizon/run
#!/usr/bin/with-contenv sh
exec s6-setuidgid php /usr/bin/php /www/artisan horizon
